# roles/diskstorage-settings/files/hdparm.service
[Unit]
Description=hdparm sleep

[Service]
Type=oneshot
DeviceAllow=/dev/sdb
ExecStart=/usr/sbin/hdparm -S 141 -y /dev/sdb

[Install]
WantedBy=multi-user.target